# packer/files/magento-salesforce-consumer.ini
[program:magento-salesforce-consumer]
user=ec2-user
autostart = true
autorestart = true
process_name = %(program_name)s_%(process_num)s
command = php /var/www/html/current/bin/magento messenger:consume async --time-limit=3600 --sleep=10
stdout_logfile = /var/www/html/current/var/log/supervisor-magento-salesforce-consumer.log
stderr_logfile = /var/www/html/current/var/log/supervisor-magento-salesforce-consumer-err.log
numprocs = 6
startretries = 3
